# vault/world/systemd/world-ingest-supervisor.service.d/throughput.conf
# Throughput tunables for the wiki ingest.
#
# Why these values:
# - BATCH=512 (from 128): fewer flush cycles -> fewer LanceDB
#   fragments -> less compaction churn. At 128 we observed ~1,650
#   fragments/hour, past the ~10k-fragment throughput-collapse
#   threshold within ~6h between compaction runs. 4x batch = 4x
#   fewer flushes. GPU embed time per batch grows linearly so this
#   isn't a GPU win, but per-flush LanceDB write overhead is
#   amortised over 4x more articles. Memory cost ~hundreds of MB
#   extra in-flight.
# - CONCURRENCY=2 (from 1): one batch can embed while another flushes
#   to LanceDB. Hides DB-write latency behind GPU work.
# - PREFETCH_QUEUE=64 (from 24): observed queue=24/24 saturation in the
#   pre-tune logs means the parser was being blocked. Larger buffer
#   lets the parser run ahead so the writer never starves.
[Service]
Environment=VAULT_INGEST_BATCH=512
Environment=VAULT_INGEST_CONCURRENCY=2
Environment=VAULT_INGEST_PREFETCH_QUEUE=64
