# node/systemd/templates/audio.service.tmpl
[Unit]
Description=LUHKAS Audio Node ({NODE_ID})
After=network-online.target sound.target {NODE_ID}-presence.service
Wants=network-online.target sound.target
{EXTRA_UNIT}
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={NODE_DIR}
Environment=LUHKAS_NODE_ID={NODE_ID}
Environment=AUDIO_STT_ENGINE=vosk
Environment=AUDIO_TTS_ENGINE=espeak
Environment=AUDIO_DISPLAY_URL=http://127.0.0.1:5005
Environment=AUDIO_PRESENCE_URL=http://127.0.0.1:5002/presence/message
{EXTRA_ENV}
ExecStart={NODE_DIR}/scripts/start_audio.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
