# node/systemd/templates/battery.service.tmpl
[Unit]
Description=Battery Node Service ({NODE_ID})
After=network-online.target
Wants=network-online.target
{EXTRA_UNIT}
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={NODE_DIR}
Environment=LUHKAS_NODE_ID={NODE_ID}
{EXTRA_ENV}
ExecStart={NODE_DIR}/scripts/start_battery.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
