# node/systemd/templates/display.service.tmpl
[Unit]
Description=LUHKAS Display Node ({NODE_ID})
After=network-online.target
Wants=network-online.target
{EXTRA_UNIT}
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={NODE_DIR}
Environment=LUHKAS_NODE_ID={NODE_ID}
Environment=AUDIO_SERVICE_URL=http://127.0.0.1:5004
{EXTRA_ENV}
ExecStart={NODE_DIR}/scripts/start_display.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
