# node/systemd/templates/rover.service.tmpl
[Unit]
Description=LUHKAS Rover Node Service ({NODE_ID})
After=network-online.target
Wants=network-online.target
{EXTRA_UNIT}
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={NODE_DIR}
Environment=LUHKAS_NODE_ID={NODE_ID}
Environment=VISION_URL=http://127.0.0.1:5000
Environment=ROBOT_API_URL=http://127.0.0.1:5001
{EXTRA_ENV}
ExecStart={NODE_DIR}/scripts/start_rover.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
