# node/systemd/templates/vision.service.tmpl
[Unit]
Description=LUHKAS Vision Service ({NODE_ID})
After=network-online.target
Wants=network-online.target
{EXTRA_UNIT}
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={NODE_DIR}
Environment=LUHKAS_NODE_ID={NODE_ID}
Environment=VAULT_CHAT_URL={VAULT_URL}
{EXTRA_ENV}
ExecStart={NODE_DIR}/scripts/start_vision.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
