# node/systemd/scout-presence.service
[Unit]
Description=Rover Unified Presence Proxy
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory={NODE_DIR}
Environment=VAULT_CHAT_URL={VAULT_URL}
Environment=VAULT_CHAT_SOURCE=scout_presence
Environment=SCOUT_PRESENCE_HOST=0.0.0.0
Environment=SCOUT_PRESENCE_PORT=5002
Environment=LUHKAS_NODE_ID={NODE_ID}
ExecStart={NODE_DIR}/scripts/start_presence.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
