# node/systemd/scout-vision.service
[Unit]
Description=Rover Vision Service
After=network-online.target scout-robot-api.service
Wants=network-online.target scout-robot-api.service
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/scout_runtime
Environment=SCOUT_TRACKING_ENABLED=1
Environment=SCOUT_TARGET_LABEL=person
Environment=SCOUT_SCORE_THRESHOLD=0.45
Environment=SCOUT_COMMAND_INTERVAL=0.12
Environment=SCOUT_ABSOLUTE_PAN_GAIN=28
Environment=SCOUT_ABSOLUTE_TILT_GAIN=22
Environment=SCOUT_ABSOLUTE_MAX_STEP=8
Environment=SCOUT_ABSOLUTE_DISTANCE_GAIN=0.6
Environment=SCOUT_ABSOLUTE_DISTANCE_MAX_MULTIPLIER=1.1
Environment=SCOUT_WHEEL_ENABLED=0
Environment=SCOUT_POSE_JOINT_THRESHOLD=0.3
Environment=SCOUT_GUARD_ALERT_URL=http://luhkas-vault.local:7000/alerts
Environment=SCOUT_GUARD_SNAPSHOT=1
Environment=SCOUT_BEHAVIOR_ENABLED=1
Environment=SCOUT_SEARCH_TIMEOUT=30
Environment=SCOUT_SEARCH_SWEEP_DURATION=3.2
Environment=SCOUT_SEARCH_SWEEP_PAN_AMOUNT=60
Environment=SCOUT_SEARCH_SCAN_PAN_PERIOD=11.0
Environment=SCOUT_SEARCH_SCAN_TILT_PERIOD=18.0
Environment=SCOUT_SEARCH_SCAN_COMMAND_INTERVAL=0.20
Environment=SCOUT_SEARCH_COMMAND_INTERVAL=0.20
Environment=SCOUT_AVOID_DURATION=3
Environment=SCOUT_EGO_MOTION_ENABLED=0
Environment=SCOUT_GYRO_PAN_SCALE=0.0
Environment=SCOUT_GYRO_TILT_SCALE=0.0
Environment=SCOUT_TELEMETRY_POLL_INTERVAL=0.2
Environment=SCOUT_VAULT_MEMORY_ENABLED=1
Environment=SCOUT_VAULT_MEMORY_URL=http://luhkas-vault.local:7000
Environment=SCOUT_TARGET_LOST_GRACE_SECONDS=1.5
Environment=SCOUT_BYTETRACKER_MATCH_THRESH=0.75
ExecStart=%h/scout_runtime/scripts/start_vision.sh
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
