# vault/world/systemd/world-compact.timer
[Unit]
Description=LUHKAS world LanceDB compaction (daily, 3:30 AM)
Documentation=man:systemd.timer(5)

[Timer]
# 3:30 AM local time daily. Coincides with low user activity (ingest is
# probably busy anyway, but optimize handles concurrent writes), well
# before the 4 AM ZIM auto-update window if that's ever configured.
OnCalendar=*-*-* 03:30:00
# If the host was off at trigger time, run as soon as possible after boot.
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
