# vault/world/systemd/world-compact.timer
[Unit]
Description=LUHKAS world LanceDB compaction (every 12h)
Documentation=man:systemd.timer(5)

[Timer]
# Every 12h: 03:30 and 15:30 local. Daily wasn't enough — at the
# observed flush rate (~1,650 fragments/hour with BATCH=128) the
# wiki_chunks table crossed the ~10k-fragment throughput-collapse
# threshold within 6h of a compaction run. With BATCH=512 the rate
# drops ~4x but we still want a safety margin, so twice-daily keeps
# fragments well under threshold throughout the cycle.
OnCalendar=*-*-* 03:30:00
OnCalendar=*-*-* 15:30:00
# If the host was off at trigger time, run as soon as possible after boot.
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
